# scripts/systemd/mediaflow-worker.service
[Unit]
Description=MediaFlow ARQ Worker (Download Queue)
After=network-online.target postgresql.service redis.service mediaflow-web.service
Wants=network-online.target
Requires=postgresql.service redis.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
User=${SERVICE_USER}
Group=${SERVICE_USER}
WorkingDirectory=${PROJECT_DIR}

EnvironmentFile=${PROJECT_DIR}/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

ExecStart=${PROJECT_DIR}/.venv/bin/python \
    -m workers.arq_worker

ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
TimeoutStartSec=30
TimeoutStopSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mediaflow-worker

# Воркер качает файлы — нужно больше памяти
LimitNOFILE=65536
MemoryMax=3G
CPUQuota=95%

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=${PROJECT_DIR}/storage
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
